# profile.d/20hook_functions.zsh
#!/bin/zsh

typeset -ga chpwd_functions
typeset -ga precmd_functions
typeset -ga preexec_functions
chpwd_functions=()
precmd_functions=()
preexec_functions=()
