# profile.d/50ansible.sh
export ANSIBLE_HOME="${XDG_CONFIG_HOME:-$HOME/.config}/ansible"
export ANSIBLE_CONFIG="${XDG_CONFIG_HOME:-$HOME/.config}/ansible.cfg"
export ANSIBLE_GALAXY_CACHE_DIR="${XDG_CACHE_HOME:-$HOME/.cache}/ansible/galaxy_cache"
export ANSIBLE_SSH_CONTROL_PATH_DIR="${XDG_RUNTIME_DIR:-/tmp}/ansible/cp"
