# profile.d/50cc.sh
export CC='ccache gcc' CXX='ccache g++'
