# profile.d/50cdpath.sh
# Search path for the cd command
CDPATH=.:~
