# profile.d/50fignore.sh
# Filename suffixes to ignore during completion
FIGNORE=\~:\$:,v:.aux:.bak:.BAK:.d:.elc:.o:.obj:.olb:.old:.orig:.rej:.so:CVS:RCS
