# profile.d/50gopath.sh
if [ -z "${GOPATH:-}" ]; then
    export GOPATH=$HOME/g
    PATH=$PATH:$GOPATH/bin
fi
