# profile.d/50gopath.zsh
path=( $path ${^${(s/:/)GOPATH}}/bin(N) )
fpath=( $fpath ${^${(s/:/)GOPATH}}/src/*/*/ghq/zsh(N) )
