# profile.d/50gzip.sh
if [ -f /etc/debian_version ]; then
    export GZIP="-9 --rsyncable"
else
    export GZIP=-9
fi
