# profile.d/50homebrew.zsh
fpath=($fpath $HOME/homebrew/share/zsh/site-functions(N))
