# profile.d/50ignoreeof.sh
# Ctrl-D でログアウトするのを抑制する。man bash 参照
IGNOREEOF=3
