# profile.d/50postgres.sh
# brew install libpq
my_append_path /opt/homebrew/opt/libpq/bin
