# profile.d/50postgres.sh
# brew install libpq
if [ -d /opt/homebrew/opt/libpq/bin ]; then
    my_append_path /opt/homebrew/opt/libpq/bin
fi
