# profile.d/50postgres.sh
# brew install libpq
if [ -d /opt/homebrew/opt/libpq/bin ]; then
    PATH="/opt/homebrew/opt/libpq/bin:$PATH"
fi
