# profile.d/50rsh.sh
export CVS_RSH=ssh
export RSYNC_RSH=$CVS_RSH
