# profile.d/50stty.sh
#stty intr "^C" kill "^U" werase "^W" susp "^Z"
stty stop "^-"
#stty erase "^?"
#stty erase "^H"
#stty -istrip
