# profile.d/50watch.zsh
#!/bin/zsh
watch=all
WATCHFMT="%D %T  $WATCHFMT"
