# profile.d/50wordchars.zsh
#!/bin/zsh
# remove '/' from WORDCHARS
WORDCHARS="${WORDCHARS:s@/@}"
