# profile.d/70npm.sh
if type npm >/dev/null 2>&1; then
  # after compinit
  eval "$(npm completion)"
fi
